# files/etc/systemd/system/bluetooth.service.d/override.conf
# https://www.redpill-linpro.com/techblog/2021/05/31/better-bluetooth-headset-audio-with-msbc.html
[Service]
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd --experimental
