# files/etc/systemd/system/rrdcached.service.d/override.conf
[Service]
# http://guide.munin-monitoring.org/en/latest/example/rrdcached/systemd.html
ExecStartPost=/bin/sleep 1 ; /bin/setfacl -m u:www-data:rw /var/run/rrdcached.sock